# artifacts/capability_tests/my-test-service.service
[Unit]
Description=AI-Generated Test Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/tmp
ExecStart=/usr/bin/python3 -m http.server 8888
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
